# cron.d/reset-dns_serial
# at noon
0 0 * * *     root  /bin/echo "0" > /home/inferno/dns/dns_serial
